# docker/supervisor/celery.ini
[program:celery]
command=celery -A celery_app worker -l INFO -B
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autorestart=true
redirect_stderr=true
stdout_logfile=/app/data/logs/celery.log
stdout_logfile_maxbytes=20MB
stdout_logfile_backups =10
loglevel=info
